# image/crontab
*/5 * * * * /bin/sh /app/download_audio.sh >> /app/config/ytdlp_audio.log 2>&1
*/7 * * * * /bin/sh /app/download_video.sh >> /app/config/ytdlp_video.log 2>&1